# MainProject/gunicorn/gunicorn.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target
[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/MainProject
ExecStart=/home/ubuntu/env/bin/gunicorn --access-logfile - --workers 3 --bind unix:/home/ubuntu/MainProject/MainProject/MainProject.sock MainProject.wsgi:application

       
[Install]
WantedBy=multi-user.target